# meta-rk3562-custom/recipes-core/display/files/rk3562-display-init.init
#!/bin/sh
### BEGIN INIT INFO
# Provides:          rk3562-display-init
# Required-Start:    $local_fs
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Prepare DRM RGB panel and /dev/fb0 for Qt linuxfb
### END INIT INFO

case "$1" in
start)
	echo "Starting rk3562-display-init..."
	/usr/sbin/rk3562-display-init.sh
	;;
stop|restart|force-reload)
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
	;;
esac

exit 0
